# ansible/playbooks/files/apcupsd_exporter.service
[Unit]
Description=APC UPS Power Prometheus Exporter

[Service]
ExecStart=/sbin/apcupsd_exporter

[Install]
WantedBy=multi-user.target
