# ansible/playbooks/files/restic/restic@.service
[Unit]
Description=restic (%I) backup

[Service]
Type=oneshot
Nice=19
IOSchedulingClass=idle
ExecStart=/etc/restic/b2_backup.sh
