# ansible/playbooks/files/cryptshot/cryptshot@.service
[Unit]
Description=cryptshot (%I) backup

[Service]
Type=oneshot
Nice=19
IOSchedulingClass=idle
ExecStart=/usr/bin/cryptshot -i %I
